# jadloop@.service
[Unit]
Description=jadloop %i

[Service]
Type=simple
ExecStart=/home/yurihs/.local/bin/jadloop notify %i

[Install]
WantedBy=default.target

